# .devcontainer/.docker/nginx-php/root/etc/service/php-fpm/run
#!/bin/sh

set -eu
exec 2>&1

if [ -f /etc/conf.d/php-fpm ]; then
    # shellcheck source=/dev/null
    . /etc/conf.d/php-fpm
fi

if [ -n "${WPVIP_XDEBUG_ENABLE}" ]; then
    sed "s/^xdebug\\.mode.*\$/xdebug.mode = ${WPVIP_XDEBUG_ENABLE}/" "${PHP_INI_DIR}/conf.d/xdebug.ini.template" > "${PHP_INI_DIR}/conf.d/xdebug.ini" || true
else
    rm -f "${PHP_INI_DIR}/conf.d/xdebug.ini" || true
fi

DAEMON=/usr/sbin/php-fpm

exec ${DAEMON} --nodaemonize
